# crawler/crontab
# Monthly crawl of all configured sources (ΦΕΚ, ΥΠΕΝ, ΤΕΕ, ...).
# Runs the 3rd of each month at 03:00 UTC — after the National Printing
# House typically finishes publishing the previous month's issues.
0 3 3 * * python -m crawler.main

# Weekly staleness sweep - flags public KB documents whose last_verified_at
# is missing or older than 6 months, for the super admin's manual review
# queue (GET /admin/stale-documents). Monday 04:00 UTC, after the sweep
# would have run against whatever the monthly crawl last touched.
0 4 * * 1 python -m crawler.staleness

# Weekly canary benchmark - 10 fixed chat questions through the real
# pipeline (5 construction, 5 tax/accounting), checking gap=false and
# citations present on the persisted chat_sessions row. Monday 05:00 UTC,
# after the staleness sweep. Alerts (failures only) go to benchmark_alerts
# and a super admin notification - see crawler/crawler/canary_benchmark.py.
0 5 * * 1 python -m crawler.canary_benchmark

# Weekly infra-health snapshot - system-wide pgvector index size (public KB
# + every company's uploads combined), classified against fixed thresholds
# and logged to infra_health_checks every run regardless of level. Monday
# 06:00 UTC, after the canary benchmark. Notifies every super admin at
# warning/critical only - infra monitoring, not billing enforcement. See
# crawler/crawler/infra_health_check.py.
0 6 * * 1 python -m crawler.infra_health_check
